# tank-data/rootfs/etc/services.d/tank-data/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the tank data service
# ==============================================================================

cd /app
exec python3 web_server.py
